# backend/healthlin-backend.service
[Unit]
Description=Healthlin Backend API
After=network.target orthanc.service
Wants=orthanc.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/healthlin/backend
ExecStart=/opt/healthlin/venv/bin/python app.py
Restart=always
RestartSec=5

Environment=ORTHANC_URL=http://localhost:8042
Environment=ORTHANC_USER=admin
Environment=ORTHANC_PASS=sua_senha
Environment=JWT_SECRET=TROQUE_ESTA_CHAVE_EM_PRODUCAO
Environment=DB_PATH=/opt/healthlin/backend/healthlin.db

[Install]
WantedBy=multi-user.target
